# archlinux/docker/install
sudo pacman -S docker
sudo systemctl enable docker.service
sudo systemctl start docker.service