# archlinux/other/etc_profile.d_bu/profile.d/ant.sh
#!/usr/bin/env sh
export ANT_HOME=/opt/apache-ant-1.9.7
export ANT=${ANT_HOME}/bin
