# archlinux/other/etc_profile.d_bu/profile.d/maven.sh
#!/usr/bin/env sh
export M2_HOME=/opt/apache-maven-3.3.9
export M2=${M2_HOME}/bin

