# archlinux/rabbitmq/install
sudo pacman -S rabbitmq
sudo systemctl enable rabbitmq.service
sudo systemctl start rabbitmq.service